# deploy/systemd/tutorsphere-web.service.template
[Unit]
Description=TutorSphere Web (Blazor)
After=network.target tutorsphere-api.service
Wants=tutorsphere-api.service

[Service]
Type=simple
WorkingDirectory=@DEPLOY_PATH@/web
ExecStart=/usr/bin/dotnet @DEPLOY_PATH@/web/@DLL_NAME@
Restart=always
RestartSec=10
User=@SERVICE_USER@
Group=@SERVICE_USER@
EnvironmentFile=@DEPLOY_PATH@/env
Environment=ASPNETCORE_URLS=http://127.0.0.1:@PORT@
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false

[Install]
WantedBy=multi-user.target
